# cuantico.service
[Unit]
Description=Cuántico — asistente de voz
# Espera a red + sonido antes de arrancar (Gemini/Deepgram/micro USB)
After=network-online.target sound.target bluetooth.target
Wants=network-online.target bluetooth.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
# NeoPixel pide /dev/mem → root obligatorio.
User=root
# Ajustado para este montaje en Raspberry Pi.
WorkingDirectory=/home/nicolas/cuantico
ExecStart=/home/nicolas/cuantico/.venv/bin/python3 /home/nicolas/cuantico/src/main.py

# Deja que el micro USB termine de enumerar
ExecStartPre=/bin/sleep 3

# Si peta (ej. wifi tarda), reintentar
Restart=on-failure
RestartSec=5

# Logs a journald (se consultan con: journalctl -u cuantico -f)
StandardOutput=journal
StandardError=journal

# Propagar señales limpiamente
KillSignal=SIGINT
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
